# jombz-web.conf
#!upstart
start on filesystem and started networking
respawn
chdir /opt/jombz-web
exec npm start >> log

