# jombz-web.conf
#!upstart
start on filesystem and started networking
respawn
chdir /opt/jombz-web
exec bash -c `source ~/.nvm/nvm.sh && npm start >> log 2>&1`

